# lib/assets/crontab
MAIL='keburgett@gmail.com'
44 10 * * 6 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  1  1 localhost:3000 1
45 10 * * 6 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  1  0 localhost:3000 1
46 10 * * 6 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  2  1 localhost:3000 2
47 10 * * 6 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  2  0 localhost:3000 2
48 10 * * 6 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  3  1 localhost:3000 3
49 10 * * 6 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  3  0 localhost:3000 3
50 10 * * 6 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  4  1 localhost:3000 4
51 10 * * 6 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  4  0 localhost:3000 4
52 10 * * 6 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  5  1 localhost:3000 5
53 10 * * 6 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  5  0 localhost:3000 5
54 10 * * 6 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  1  1 localhost:3000 6
55 10 * * 6 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  1  0 localhost:3000 6
56 10 * * 6 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  2  1 localhost:3000 7
57 10 * * 6 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  2  0 localhost:3000 7
58 10 * * 6 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  3  1 localhost:3000 8
59 10 * * 6 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  3  0 localhost:3000 8
00 11 * * 6 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  4  1 localhost:3000 9
01 11 * * 6 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  4  0 localhost:3000 9
02 11 * * 6 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  5  1 localhost:3000 10
03 11 * * 6 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  5  0 localhost:3000 10
