# lib/assets/crontab
MAIL='keburgett@gmail.com'
46 10 * * 5 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  1  1 localhost:3000 1
47 10 * * 5 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  1  0 localhost:3000 1
48 10 * * 5 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  2  1 localhost:3000 2
49 10 * * 5 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  2  0 localhost:3000 2
50 10 * * 5 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  3  1 localhost:3000 3
51 10 * * 5 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  3  0 localhost:3000 3
52 10 * * 5 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  4  1 localhost:3000 4
53 10 * * 5 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  4  0 localhost:3000 4
54 10 * * 5 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  5  1 localhost:3000 5
55 10 * * 5 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  5  0 localhost:3000 5
56 10 * * 5 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  1  1 localhost:3000 6
57 10 * * 5 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  1  0 localhost:3000 6
58 10 * * 5 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  2  1 localhost:3000 7
59 10 * * 5 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  2  0 localhost:3000 7
00 11 * * 5 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  3  1 localhost:3000 8
01 11 * * 5 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  3  0 localhost:3000 8
02 11 * * 5 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  4  1 localhost:3000 9
03 11 * * 5 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  4  0 localhost:3000 9
04 11 * * 5 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  5  1 localhost:3000 10
05 11 * * 5 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  5  0 localhost:3000 10
